# xdcterm.service
[Unit]
Description=xdcterm Delta Chat bot
After=network-online.target

[Service]
Type=simple
Environment=PATH=/home/user/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=TERM=xterm-256color
WorkingDirectory=/home/user/python/xdcterm
ExecStart=uv run --directory backend bot.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
