# xdcterm.service
[Unit]
Description=xdcterm Delta Chat bot
After=network-online.target

[Service]
Type=simple
Environment=TERM=xterm-256color
WorkingDirectory=/home/user/python/xdcterm
ExecStart=/home/user/.local/bin/uv run --directory backend bot.py serve
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
