# xdcterm.service
[Unit]
Description=xdcterm Delta Chat bot
After=network-online.target

[Service]
Type=simple
Environment=PATH=/home/user/.local/bin:/usr/local/bin:/usr/bin:/bin
WorkingDirectory=/home/user/python/xdcterm/backend
ExecStart=/home/user/.local/bin/uv run bot.py serve
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
